# systemd/snowflake_update.service
[Unit]
Description=Snowflake Proxy Update start

[Service]
ExecStart=/opt/snowflake/update_snowflake.sh
