# systemd/snowflake.service
[Unit]
Description=Snowflak Proxy
After=network-online.target

[Service]
Type=simple
ExecStart=/opt/proxy
Restart=always
User=snowflake
Group=snowflake

[Install]
WantedBy=multi-user.target
